# com.mergepdf.folder-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mergepdf.folder-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/path/to/mergepdf/folder_monitor.py</string>
        <string>--watch-dir</string>
        <string>/Users/YOURUSERNAME/Downloads</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/path/to/mergepdf/logs/folder_monitor_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/path/to/mergepdf/logs/folder_monitor_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/path/to/mergepdf</string>
</dict>
</plist>
